# launchctl/com.music.genreautoupdater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Genre Updater LaunchAgent Template

  This is a TEMPLATE. The actual plist is generated by install.sh
  with correct paths for your system.

  Location after install: ~/Library/LaunchAgents/com.music.genreautoupdater.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.music.genreautoupdater</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Path to wrapper script (generated by install.sh) -->
        <string>$HOME/Library/Application Support/GenreUpdater/bin/run-daemon.sh</string>
    </array>

    <!-- Trigger on Music Library changes -->
    <key>WatchPaths</key>
    <array>
        <string>$HOME/Music/Music/Music Library.musiclibrary</string>
    </array>

    <key>StandardOutPath</key>
    <string>$HOME/Library/Application Support/GenreUpdater/logs/launchctl-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>$HOME/Library/Application Support/GenreUpdater/logs/launchctl-stderr.log</string>

    <!-- Minimum 5 minutes between triggers (launchctl level) -->
    <key>ThrottleInterval</key>
    <integer>300</integer>

    <!-- Don't run on login, only on WatchPaths trigger -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Not a persistent daemon, exits after each run -->
    <key>KeepAlive</key>
    <false/>

    <!-- Low CPU priority -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- 4 hour max runtime -->
    <key>ExitTimeOut</key>
    <integer>14400</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
